# supervisord-gammu.conf
[program:gammu]
command=/start-gammu.sh
numprocs=1
autostart=true
autorestart=false
